# overlay/etc/profile.d/k3os-system.sh
alias k3os-system='kubectl -n k3os-system'